# status/server.conf
[supervisord]
nodaemon=true

[program:server]
command=/server.py
