# status/status_check.conf
[supervisord]
nodaemon=true

[program:status_check]
command=/status_check.py
