# nginx-reverse-proxy/supervisor.conf
[supervisord]
nodaemon=true

[program:update_config]
command=/update_config.py

